# openclaw-watchdog.service
[Unit]
Description=OpenClaw lightweight watchdog

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/openclaw-watchdog.sh
